# packaging/debian/postinst
#!/usr/bin/env bash
set -e

# Aviso não-fatal: a versão mínima de Node.js não é imposta pelo apt
# (para não brigar com o pacote `nodejs` de cada distro/PPA), então checamos
# manualmente e só avisamos — não abortamos a instalação por causa disso.
NODE_BIN="$(command -v node || true)"
if [ -n "$NODE_BIN" ]; then
  NODE_MAJOR="$("$NODE_BIN" -p 'process.versions.node.split(".")[0]' 2>/dev/null || echo 0)"
  if [ "${NODE_MAJOR:-0}" -lt 20 ] 2>/dev/null; then
    echo "Aviso: Héstia Console requer Node.js 20+; detectado $("$NODE_BIN" -v)." >&2
  fi
else
  echo "Aviso: binário 'node' não encontrado no PATH; o serviço não vai conseguir iniciar." >&2
fi

# Criar /etc/default/hestia-console se não existir
if [ ! -f /etc/default/hestia-console ]; then
  echo "Criando template de configuração em /etc/default/hestia-console"
  cat << 'EOF' > /etc/default/hestia-console
# Héstia Console Configuration
# Future Station URL placeholder. Not used yet.
# HESTIA_STATION_BASE_URL=https://station.example.ts.net

# Allowed Host headers (comma-separated, exact values, no wildcards)
# HESTIA_ALLOWED_HOSTS=hestia.exemplo.ts.net,hestia.exemplo.ts.net:443
EOF
fi

# Criar grupo de sistema se não existir
if ! getent group hestia-console >/dev/null; then
  echo "Criando grupo de sistema hestia-console..."
  addgroup --system hestia-console
fi

# Criar usuário de sistema se não existir
if ! getent passwd hestia-console >/dev/null; then
  echo "Criando usuário de sistema hestia-console..."
  adduser --system --ingroup hestia-console --no-create-home --home /opt/hestia-console --shell /usr/sbin/nologin hestia-console
fi

systemctl daemon-reload
systemctl enable hestia-console.service
systemctl restart hestia-console.service

gtk-update-icon-cache -q /usr/share/icons/hicolor || true
update-desktop-database /usr/share/applications || true

echo "Héstia Console instalada."
echo "Abra pelo menu ou acesse http://localhost:4517"
