# packaging/debian/postinst
#!/usr/bin/env bash
set -e

# Aviso não-fatal: a versão mínima de Node.js não é imposta pelo apt
# (para não brigar com o pacote `nodejs` de cada distro/PPA), então checamos
# manualmente e só avisamos — não abortamos a instalação por causa disso.
NODE_BIN="$(command -v node || true)"
if [ -n "$NODE_BIN" ]; then
  NODE_MAJOR="$("$NODE_BIN" -p 'process.versions.node.split(".")[0]' 2>/dev/null || echo 0)"
  if [ "${NODE_MAJOR:-0}" -lt 20 ] 2>/dev/null; then
    echo "Aviso: Héstia Console requer Node.js 20+; detectado $("$NODE_BIN" -v)." >&2
  fi
else
  echo "Aviso: binário 'node' não encontrado no PATH; o serviço não vai conseguir iniciar." >&2
fi

systemctl daemon-reload
systemctl enable hestia-console.service
systemctl restart hestia-console.service

gtk-update-icon-cache -q /usr/share/icons/hicolor || true
update-desktop-database /usr/share/applications || true

# O serviço roda com DynamicUser=yes (sem privilégio de root) e escreve de verdade em /KALINE
# (organizer). Se /KALINE já existir neste momento, libera escrita pro grupo do serviço; se
# ainda não existir, isso é só um aviso — repita manualmente depois de criar /KALINE.
if [ -d /KALINE ]; then
  chgrp hestia-console /KALINE 2>/dev/null && chmod g+rwx /KALINE 2>/dev/null \
    || echo "Aviso: não consegui ajustar permissões de /KALINE automaticamente. Rode manualmente: sudo chgrp hestia-console /KALINE && sudo chmod g+rwx /KALINE" >&2
else
  echo "Aviso: /KALINE ainda não existe. Depois de criá-la, rode: sudo chgrp hestia-console /KALINE && sudo chmod g+rwx /KALINE" >&2
fi

echo "Héstia Console instalada."
echo "Abra pelo menu ou acesse http://localhost:4517"
