# packaging/debian/prerm
#!/usr/bin/env bash
set -e

systemctl stop hestia-console.service || true
systemctl disable hestia-console.service || true
