# packaging/hestia-console.service.in
[Unit]
Description=Héstia Console (rodando direto do checkout do git em __WORKDIR__)
After=network.target

[Service]
Type=simple
WorkingDirectory=__WORKDIR__
ExecStart=/usr/bin/env node __WORKDIR__/hestia.js
Restart=on-failure
RestartSec=3
StateDirectory=hestia-console
LogsDirectory=hestia-console
Environment=HESTIA_HOST=127.0.0.1
Environment=HESTIA_PORT=4517
Environment=NODE_ENV=production
EnvironmentFile=-/etc/default/hestia-console

User=hestia-console
Group=hestia-console
NoNewPrivileges=yes
ProtectSystem=full

[Install]
WantedBy=multi-user.target
