# packaging/hestia-station-agent.service.in
[Unit]
Description=Héstia Station Agent
After=network.target

[Service]
Type=simple
WorkingDirectory=__WORKDIR__
ExecStart=__NODE_BIN__ __WORKDIR__/station.js
User=__SERVICE_USER__
Group=__SERVICE_GROUP__
Environment=NODE_ENV=production
EnvironmentFile=-/etc/default/hestia-station-agent
Restart=on-failure
RestartSec=3
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full

[Install]
WantedBy=multi-user.target
