# packaging/hestia-console.service.in
[Unit]
Description=Héstia Console (rodando direto do checkout do git em __WORKDIR__)
After=network.target

[Service]
Type=simple
WorkingDirectory=__WORKDIR__
ExecStart=/usr/bin/env node __WORKDIR__/hestia.js
Restart=on-failure
RestartSec=3
StateDirectory=hestia-console
LogsDirectory=hestia-console
Environment=HESTIA_HOST=127.0.0.1
Environment=HESTIA_PORT=4517
Environment=NODE_ENV=production
EnvironmentFile=-/etc/default/hestia-console

# Roda sem privilégio de root — ver comentário equivalente em hestia-console.service.
# IMPORTANTE pra este caminho (checkout de git, não /opt): __WORKDIR__ precisa ser legível
# pelo UID efêmero do DynamicUser, o que normalmente NÃO é o caso se o checkout estiver dentro
# de um $HOME com permissão 700/750 (padrão). scripts/install.sh avisa sobre isso; o jeito mais
# simples é `chmod o+rX` no diretório do checkout e nos pais até a raiz, ou clonar em um path
# fora do home (ex.: /srv/hestia-console).
DynamicUser=yes
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=-/KALINE

[Install]
WantedBy=multi-user.target
